# bio4d/blobio.bio4d.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
 <key>Label</key>
 <string>blobio.bio4d</string>

 <key>EnvironmentVariables</key>
 <dict>
  <key>PATH</key>
  <string>/usr/local/jmscott/bin:/usr/bin:/bin</string>

  <key>BLOBIO_ROOT</key>
  <string>/usr/local/blobio</string>

  <key>BLOBIO_BIO4D_PORT</key>
  <string>1797</string>

  <!-- 
    Write round robin stats every 60 seconds to $BLOBIO_ROOT/log/bio4d-stat.log.
    0 disables logging of stat samples.
  -->
  <key>BLOBIO_BIO4D_STAT_SAMPLE_DURATION</key>
  <string>60</string>
 </dict>

 <key>KeepAlive</key>
 <true/>

 <key>ProcessType</key>
 <string>Background</string>

 <key>AbandonProcessGroup</key>
 <false/>
 <key>ProgramArguments</key>
 <array>
  <string>/usr/local/jmscott/sbin/launchd-log</string>
  <string>/usr/local/blobio/sbin/launchd-bio4d</string>
 </array>

 <key>RunAtLoad</key>
 <true/>

 <key>UserName</key>
 <string>blobio</string>

 <key>GroupName</key>
 <string>staff</string>

 <key>WorkingDirectory</key>
 <string>/usr/local/blobio</string>

 <key>StandardErrorPath</key>
 <string>/Users/blobio/log/launchd-bio4d.log</string>

 <key>StandardOutPath</key>
 <string>/Users/blobio/log/launchd-bio4d.log</string>

</dict>
</plist>
