# bio4d/blobio.bio4d.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<!--
  Synopsis:
  	Example of Mac launchd configurations for bio4d daemon.
  Usage:
  	cat blobio.bio4d.plist.example				|
  	sed 's@{BLOBIO_ROOT}@/path/to/blobio@g' 		|
	sed 's/{BLOBIO_USER}/jmscott/g'				|
	sed 's/{BLOBIO_GROUP}/staff/g'				|
	sed 's/{BLOBIO_LAUNCHD_NAME}/jmsdesk.blobio.bio4d/g'	\
		>/Library/LaunchDaemons/{BLOBIO_LAUNCHD_NAME}.plist
  See:
  	Requires jmscott work tools on github.

  		https://github.com/jmscott/work
  Note:
  	Ought to write script sed-blobio.bio4d.plist, to paramterize above
	sed script.
-->
<plist version="1.0">
<dict>
 <key>Label</key>
 <string>{BLOBIO_LAUNCHD_NAME}</string>

 <!--
	All environment vars defined in $BLOBIO_ROOT/etc/profile,
	loaded by script $BLOBIO_ROOT/sbin/launchd-bio4d.
 -->
 <key>EnvironmentVariables</key>
 <dict>
  
  <key>BLOBIO_ROOT</key>
  <string>{BLOBIO_ROOT}</string>

 </dict>

 <key>KeepAlive</key>
 <true/>

 <key>ProcessType</key>
 <string>Background</string>

 <key>AbandonProcessGroup</key>
 <false/>

 <key>ProgramArguments</key>
 <array>
  <string>/usr/local/jmscott/sbin/launchd-log</string>
  <string>{BLOBIO_ROOT}/sbin/launchd-bio4d</string>
 </array>

 <key>RunAtLoad</key>
 <true/>

 <key>UserName</key>
 <string>{BLOBIO_USER}</string>

 <key>GroupName</key>
 <string>{BLOBIO_GROUP}</string>

 <key>WorkingDirectory</key>
 <string>{BLOBIO_ROOT}</string>

 <key>StandardErrorPath</key>
 <string>/Users/{BLOBIO_USER}/log/launchd-jmsdesk-bio4d.log</string>

 <key>StandardOutPath</key>
 <string>/Users/{BLOBIO_USER}/log/launchd-jmsdesk-bio4d.log</string>
 
 <key>Disabled</key>
 <false/>

</dict>
</plist>
